# deploy/common/systemd/kiosk-shell.service
[Unit]
Description=Event Kiosk Electron Shell (X11/Wayland session)
After=kiosk-web.service network.target
Requires=kiosk-web.service

[Service]
Type=simple
User=kiosk
Environment=KIOSK_ROOT=/opt/kiosk
Environment=KIOSK_DATA_DIR=/var/lib/kiosk
Environment=KIOSK_URL=http://localhost:3000/kiosk
Environment=DISPLAY=:0
EnvironmentFile=-/var/lib/kiosk/display.env
WorkingDirectory=/opt/kiosk/shell
ExecStartPre=/bin/sleep 5
ExecStartPre=/opt/kiosk/bin/set-display-rotation.sh --apply-x11
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
