# deploy/fedora-atomic/quadlet/kiosk-web.container
[Unit]
Description=Event Kiosk Web Application
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=event-kiosk-web
Image=ghcr.io/bytelake/event-kiosk-web:latest
PublishPort=3000:3000
Volume=/var/lib/kiosk:/app/data:Z
Volume=/var/lib/kiosk/.env:/app/.env:ro,Z
Environment=NODE_ENV=production
Environment=HOSTNAME=0.0.0.0
Environment=PORT=3000
Environment=DATABASE_URL=file:/app/data/kiosk.db
Environment=UPLOADS_DIR=/app/data/uploads
AutoUpdate=registry

[Service]
Restart=always
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
